# agent/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create printmaster user and group if they don't exist
        if ! getent group printmaster >/dev/null 2>&1; then
            groupadd --system printmaster
        fi
        if ! getent passwd printmaster >/dev/null 2>&1; then
            useradd --system --gid printmaster --home-dir /var/lib/printmaster \
                --shell /usr/sbin/nologin --comment "PrintMaster Agent" printmaster
        fi

        # Set ownership of data directories
        chown -R printmaster:printmaster /var/lib/printmaster
        chown -R printmaster:printmaster /var/log/printmaster
        chmod 750 /var/lib/printmaster
        chmod 750 /var/log/printmaster

        # Create default config from example if not exists
        if [ ! -f /etc/printmaster/agent.toml ]; then
            cp /etc/printmaster/agent.toml.example /etc/printmaster/agent.toml
            chown root:printmaster /etc/printmaster/agent.toml
            chmod 640 /etc/printmaster/agent.toml
        fi

        # Ensure sudoers file has correct permissions
        if [ -f /etc/sudoers.d/printmaster-agent ]; then
            chmod 440 /etc/sudoers.d/printmaster-agent
            chown root:root /etc/sudoers.d/printmaster-agent
        fi

        # Reload systemd, enable and start/restart service
        systemctl daemon-reload
        systemctl enable printmaster-agent.service || true
        # Use restart to handle both fresh installs and upgrades
        # (start only works if service is not running at all)
        systemctl restart printmaster-agent.service || true
        
        echo "PrintMaster Agent installed and running at http://localhost:8080"
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
