# agent/debian/prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade|deconfigure)
        # Stop and disable the service before removal
        if systemctl is-active --quiet printmaster-agent.service; then
            systemctl stop printmaster-agent.service || true
        fi
        if systemctl is-enabled --quiet printmaster-agent.service 2>/dev/null; then
            systemctl disable printmaster-agent.service || true
        fi
        ;;
    failed-upgrade)
        ;;
    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
